# data/injects.json
[
  {
    "id": "inj1",
    "order": 1,
    "time": "05:00",
    "title": "A Door That Shouldn't Be Open",
    "phase": "Detection",
    "nist": "Detection & Analysis",
    "description": "A security alert goes off. An automated account, the kind no person ever logs into directly, has just logged in from a finance employee's computer to the system that controls company-wide logins. That's unusual: when an automatic account suddenly starts acting like a person, it usually means someone else is controlling it.",
    "evidence": "SIEM: 4624 logon type 3, service account svc_backup, source FIN-WS-042, 09:19:03.",
    "choices": [
      {
        "id": "correct",
        "label": "Disconnect the affected computer from the network (without turning it off), lock the automatic account so it can't be used again, check the logs to see if the intruder reached other computers, and open a formal security case.",
        "type": "correct"
      },
      {
        "id": "wrong",
        "label": "Dismiss the alert as a glitch, assuming automatic accounts sometimes act strange, and take no further action.",
        "type": "wrong"
      },
      {
        "id": "neutral",
        "label": "Mark the alert as suspicious and ask someone to look into it later in the shift, without disconnecting anything yet.",
        "type": "neutral"
      }
    ],
    "facilitatorNotes": "Watch whether the team treats this as noise or investigates. Good teams pivot on the account and host immediately.",
    "hint": "A service account interactively used from a user workstation is rarely legitimate.",
    "mitre": "TA0001 Initial Access / T1078 Valid Accounts"
  },
  {
    "id": "inj2",
    "order": 2,
    "time": "10:00",
    "title": "Excel Doesn't Do That",
    "phase": "Detection",
    "nist": "Detection & Analysis",
    "description": "Minutes later, the security software raises a serious alarm on the same computer: opening a spreadsheet has somehow triggered a hidden command window to run scrambled, disguised instructions in the background. Spreadsheets don't do that on their own; this is a sign of an attack already underway.",
    "evidence": "EDR process tree: EXCEL.EXE \u2192 powershell.exe -enc <redacted>. Network beacon to 185.xx.xx.xx every 60s.",
    "choices": [
      {
        "id": "correct",
        "label": "Treat this as a confirmed attack, automatically isolate the computer using the security software, decode the hidden instructions to understand what they were trying to do, and check the rest of the network for the same pattern.",
        "type": "correct"
      },
      {
        "id": "wrong",
        "label": "Dismiss it as a harmless glitch caused by normal spreadsheet features, and let the employee keep working on the computer.",
        "type": "wrong"
      },
      {
        "id": "neutral",
        "label": "Note the alert and run a basic virus scan on the computer, but wait for more proof before disconnecting it.",
        "type": "neutral"
      }
    ],
    "facilitatorNotes": "Assess whether they preserve evidence before isolating. Both matter, sequencing is a discussion point.",
    "hint": "Host isolation via EDR keeps forensic data intact while cutting network access.",
    "mitre": "TA0002 Execution / T1059.001 PowerShell"
  },
  {
    "id": "inj3",
    "order": 3,
    "time": "15:00",
    "title": "It's Not Just One Machine Anymore",
    "phase": "Analysis",
    "nist": "Detection & Analysis",
    "description": "The investigation traces the problem back to an email attachment made to look like a routine invoice. Someone opened it, not realizing it quietly installed malicious software. That same malicious software is now showing up on two more computers. What started as one odd login has become a foothold spreading across the network.",
    "evidence": "Email gateway log of the phishing message; sandbox detonation report; beacon seen on FIN-WS-042, FIN-WS-018, HR-WS-007.",
    "choices": [
      {
        "id": "correct",
        "label": "Officially declare a security incident, disconnect all three affected computers from the network, block the sender and the malicious file company-wide, start figuring out how far the attacker has reached, and alert the person in charge of the response.",
        "type": "correct"
      },
      {
        "id": "neutral",
        "label": "Disconnect only the first affected computer and make a note to check the other two later, without officially declaring an incident yet.",
        "type": "neutral"
      },
      {
        "id": "wrong",
        "label": "Treat each affected computer as a separate, unrelated problem, without connecting the dots or telling leadership.",
        "type": "wrong"
      }
    ],
    "facilitatorNotes": "This is the key declaration moment. Note the time; success criteria target declaration by 15:00.",
    "hint": "Multiple beaconing hosts = confirmed intrusion. Declare and coordinate now.",
    "mitre": "TA0001 Initial Access / T1566.001 Spearphishing Attachment"
  },
  {
    "id": "inj4",
    "order": 4,
    "time": "20:00",
    "title": "Finance Goes Dark",
    "phase": "Containment",
    "nist": "Containment, Eradication & Recovery",
    "description": "The server that stores the finance team's files suddenly works extremely hard, then goes silent. All of its shared folders disappear. Across the finance department, screens freeze and nobody can open a single document. On the surface it looks like a normal tech outage. In reality, it's the same attack wearing a different disguise.",
    "evidence": "Monitoring: FIN-SRV-01 SMB service down at 09:33. Scheduled task SysUpdate created 4 minutes earlier by svc_backup.",
    "choices": [
      {
        "id": "neutral",
        "label": "Restart the server right away to get things working again, without checking with the security team or saving any evidence first.",
        "type": "neutral"
      },
      {
        "id": "wrong",
        "label": "Treat it as an everyday IT outage, follow the normal restart-and-update troubleshooting steps, and don't connect it to the security incident already underway.",
        "type": "wrong"
      },
      {
        "id": "correct",
        "label": "Recognize this outage is connected to the ongoing attack rather than a random glitch, save all logs and evidence before restarting anything, and have IT and security check together whether files are actively being locked up.",
        "type": "correct"
      }
    ],
    "facilitatorNotes": "Tension between isolating (stops spread) and business impact (halts finance). Look for risk-based reasoning.",
    "hint": "A suspicious scheduled task by a compromised service account is a common ransomware staging step.",
    "mitre": "TA0040 Impact / T1053.005 Scheduled Task"
  },
  {
    "id": "inj5",
    "order": 5,
    "time": "25:00",
    "title": "The Phones Start Ringing",
    "phase": "Communication",
    "nist": "Containment, Eradication & Recovery",
    "description": "Word is spreading outside the building before anyone inside has agreed what to say. The customer support desk is flooded with complaints. Clients can't get their usual reports and want answers immediately. The staff who talk to clients don't know what to tell them, and are asking the response team for guidance before they say the wrong thing.",
    "evidence": "Service desk queue: 37 tickets in 10 minutes referencing reports unavailable. Two clients cite contractual SLAs.",
    "choices": [
      {
        "id": "wrong",
        "label": "Let client-facing staff tell customers the company has been hacked and guess about stolen data, before anything has actually been confirmed.",
        "type": "wrong"
      },
      {
        "id": "correct",
        "label": "Give client-facing staff one approved, careful statement to use, and have all communication with clients go through a single coordinated source, without guessing at causes.",
        "type": "correct"
      },
      {
        "id": "neutral",
        "label": "Let each staff member explain the situation to clients in their own words, which leads to inconsistent and confusing messages.",
        "type": "neutral"
      }
    ],
    "facilitatorNotes": "Does the team activate communications early, or stay heads-down on tech? Parallel workstreams are the goal.",
    "hint": "External pressure will grow, so get an approved holding statement ready before you need it."
  },
  {
    "id": "inj6",
    "order": 6,
    "time": "30:00",
    "title": "The Note",
    "phase": "Analysis",
    "nist": "Containment, Eradication & Recovery",
    "description": "Then it appears, on several computers at once, unmistakable: a ransom note. The attackers are demanding a large cryptocurrency payment within 72 hours. They claim to have already stolen a large amount of client data and are threatening to release it publicly if they aren't paid. The clock is now running, whether the company is ready or not.",
    "evidence": "Ransom note text; TOR negotiation URL; sample of allegedly exfiltrated filenames matching client folders.",
    "choices": [
      {
        "id": "correct",
        "label": "Save the ransom note as evidence, avoid contacting the attackers directly, and immediately tell the legal team, senior leadership, the cyber insurance provider, and law enforcement.",
        "type": "correct"
      },
      {
        "id": "neutral",
        "label": "Share the note internally and casually discuss whether to pay, without yet contacting law enforcement or the insurance provider.",
        "type": "neutral"
      },
      {
        "id": "wrong",
        "label": "Have one staff member open the attackers' payment page and start negotiating on their own, without approval from leadership or the legal team.",
        "type": "wrong"
      }
    ],
    "facilitatorNotes": "Surface the ransom-payment question and who has authority. Sanctions/insurer implications should come up.",
    "hint": "Exfiltration claims change notification obligations, so validate, but assume data loss is possible.",
    "mitre": "TA0040 Impact / T1486 Data Encrypted for Impact"
  },
  {
    "id": "inj7",
    "order": 7,
    "time": "35:00",
    "title": "A Reporter Is On the Line",
    "phase": "Communication",
    "nist": "Post-Incident Activity",
    "description": "Somehow, word has already reached a journalist. A reporter emails the company's general inbox asking them to confirm the outage affecting clients, and says they're on a deadline in two hours. Whatever the company says next, or doesn't say, will end up in print.",
    "evidence": "Reporter email with newsroom signature; a client has apparently posted about the outage on social media.",
    "choices": [
      {
        "id": "wrong",
        "label": "Have an untrained employee confirm, on the record, that the company was hacked and hit with ransomware.",
        "type": "wrong"
      },
      {
        "id": "neutral",
        "label": "Let the deadline pass without responding, letting the story run without any input from the company.",
        "type": "neutral"
      },
      {
        "id": "correct",
        "label": "Have the communications team issue a short, careful statement acknowledging a service disruption that's under investigation, without confirming any details, after checking it with the legal team first.",
        "type": "correct"
      }
    ],
    "facilitatorNotes": "Check for message discipline; no one should freelance to the press. Legal review of external statements.",
    "hint": "Acknowledge awareness and commitment to clients without confirming unverified specifics.",
    "mitre": "TA0040 Impact / T1486 Data Encrypted for Impact"
  },
  {
    "id": "inj8",
    "order": 8,
    "time": "40:00",
    "title": "The Legal Clock Starts Ticking",
    "phase": "Communication",
    "nist": "Post-Incident Activity",
    "description": "While the technical team is still working the problem, the legal team delivers unwelcome news: the stolen data likely includes personal information about customers in Europe and financial records for customers in the US. That fact alone triggers a whole set of legal notification requirements, each with its own deadline, and none of them will wait for the investigation to finish.",
    "evidence": "Sample exfiltrated filenames include EU client folders; data map indicates PII and financial account data.",
    "choices": [
      {
        "id": "neutral",
        "label": "Have the legal team keep this in mind, but wait until they're completely certain of the scope before preparing any notifications.",
        "type": "neutral"
      },
      {
        "id": "correct",
        "label": "Have the legal team start tracking notification deadlines right away, work out which data and which people are affected while the technical investigation continues, and start drafting the required notifications with outside lawyers.",
        "type": "correct"
      },
      {
        "id": "wrong",
        "label": "Leave the legal deadlines unassigned, assuming notifications can wait until the technical problem is completely fixed.",
        "type": "wrong"
      }
    ],
    "facilitatorNotes": "Do they know their obligations and timelines? Note whether cyber insurance is engaged early (often required).",
    "hint": "The GDPR 72-hour clock generally starts at awareness of a personal-data breach.",
    "mitre": "TA0040 Impact / T1486 Data Encrypted for Impact"
  },
  {
    "id": "inj9",
    "order": 9,
    "time": "45:00",
    "title": "The Room Goes Quiet",
    "phase": "Communication",
    "nist": "Post-Incident Activity",
    "description": "Senior leaders gather. On the screen behind them, a countdown ticks toward the attackers' deadline. Three decisions can no longer be put off: whether to pay, how to recover, and what (if anything) to tell the public. Everyone in the room knows this decision will define how the company is remembered.",
    "evidence": "Concise IC brief: scope, business impact (ZMW 2,000,000/hr), backup status pending, legal/insurer position on payment.",
    "choices": [
      {
        "id": "wrong",
        "label": "Have one senior leader decide alone to pay the ransom, without consulting the legal team, the insurance provider, or checking whether paying could be illegal.",
        "type": "wrong"
      },
      {
        "id": "correct",
        "label": "Have leadership make a documented decision only after consulting the legal team, the cyber insurance provider, and law enforcement, and write down the reasoning behind the choice.",
        "type": "correct"
      },
      {
        "id": "neutral",
        "label": "Keep putting off the decision, with no firm answer given, while the countdown keeps running.",
        "type": "neutral"
      }
    ],
    "facilitatorNotes": "The IC brief should be decision-oriented and brief. Executives set risk tolerance; avoid rushing to pay.",
    "hint": "Paying does not guarantee recovery and may carry sanctions risk, so weigh it against tested backups.",
    "mitre": "TA0040 Impact / T1486 Data Encrypted for Impact"
  },
  {
    "id": "inj10",
    "order": 10,
    "time": "50:00",
    "title": "What's Actually Recoverable",
    "phase": "Recovery",
    "nist": "Containment, Eradication & Recovery",
    "description": "Amid the chaos, some good news arrives, with a catch. The IT team reports that a secure, untouched backup of the finance server from last night appears intact. Recovery is possible. But restoring it will take 6 to 8 hours, and worse, two other affected servers have no recent backup at all.",
    "evidence": "Backup console: FIN-SRV-01 last good immutable snapshot 02:00 today; APP-SRV-03 last backup 9 days ago; DEV-SRV-02 no backups.",
    "choices": [
      {
        "id": "correct",
        "label": "Check that the backup is clean and trustworthy before restoring it, prioritize which systems to recover first based on how critical they are to the business, and flag the two servers with no backup as an urgent problem needing an immediate fix.",
        "type": "correct"
      },
      {
        "id": "wrong",
        "label": "Restore every system immediately using whatever backups exist, without first checking that they're free of the attacker's software.",
        "type": "wrong"
      },
      {
        "id": "neutral",
        "label": "Restore the finance server without first checking that the backup is clean, with no plan yet for the two servers that have no backup.",
        "type": "neutral"
      }
    ],
    "facilitatorNotes": "Look for verification before restore and eradication-before-recovery discipline. Surface backup gaps as findings.",
    "hint": "Never restore into a still-compromised environment; eradicate first, verify backups are clean.",
    "mitre": "TA0040 Impact / T1486 Data Encrypted for Impact"
  },
  {
    "id": "inj11",
    "order": 11,
    "time": "55:00",
    "title": "Catching Breath",
    "phase": "Recovery",
    "nist": "Post-Incident Activity",
    "description": "The facilitator calls for a pause. The immediate crisis is over, but the real reflection is just beginning. Each team involved (security, IT, legal, communications, leadership) is asked the same question: what worked well, and what weakness did you just discover about how prepared you really were?",
    "evidence": "Decision log to date; observer notes; open action items.",
    "choices": [
      {
        "id": "correct",
        "label": "Have each team share one specific thing that worked well and one specific weakness they found, with someone assigned to follow up in the final report.",
        "type": "correct"
      },
      {
        "id": "neutral",
        "label": "Give vague feedback, like saying communication could have been better, without any specifics or anyone assigned to fix it.",
        "type": "neutral"
      },
      {
        "id": "wrong",
        "label": "Have teams get defensive and downplay their mistakes, with nothing written down at all.",
        "type": "wrong"
      }
    ],
    "facilitatorNotes": "Keep it blameless and specific. Convert observations into concrete, ownable actions for the report.",
    "hint": "Focus on systemic fixes (process, tooling, training) rather than individual mistakes.",
    "mitre": "TA0040 Impact / T1486 Data Encrypted for Impact"
  },
  {
    "id": "inj12",
    "order": 12,
    "time": "60:00",
    "title": "Closing the Book",
    "phase": "Recovery",
    "nist": "Post-Incident Activity",
    "description": "The exercise wraps up. The facilitator has one last job: make sure the scoring is fair, make sure the record of what happened is accurate, and produce the final report that will decide whether anything learned today actually leads to change.",
    "evidence": "Completed scorecard; full decision log; consolidated recommendations.",
    "choices": [
      {
        "id": "wrong",
        "label": "Rush or skip the scoring after everyone leaves, leave the record of decisions incomplete, and never produce or share a final report.",
        "type": "wrong"
      },
      {
        "id": "neutral",
        "label": "Finish the scoring quickly with limited input, and produce a final report with vague recommendations that nobody is assigned to follow up on.",
        "type": "neutral"
      },
      {
        "id": "correct",
        "label": "Finalize the scoring with input from participants, double-check the decision record against the full timeline, and produce a final report with clear follow-up actions, owners, and deadlines.",
        "type": "correct"
      }
    ],
    "facilitatorNotes": "Thank participants, reinforce the no-fault ethos, and set expectations for the written report and follow-ups.",
    "hint": "A tabletop only pays off if actions are tracked to closure, so assign owners and due dates.",
    "mitre": "TA0040 Impact / T1486 Data Encrypted for Impact"
  }
]